# etc/udev/rules.d/60-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add", RUN+="/bin/chgrp staff $sys$devpath/brightness", RUN+="/bin/chmod g+w $sys$devpath/brightness"
